# server/bin/vaerhona-staging
#! /bin/sh
### BEGIN INIT INFO
# Provides: vaerhona-staging
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Værhøna staging enviroment
# Description: This file starts and stops Værhøna staging enviroment
# 
### END INIT INFO

APP_DIR=/var/www/html/vaerhona

case "$1" in
 start)
   PORT=8001 NODE_ENV=production AWS_PROFILE=staging MONGO_URL=mongodb://localhost:27017/vaerhona-staging forever start --sourceDir=$APP_DIR -w --watchDirectory $APP_DIR/dist --silent --uid "vaerhona-staging" --append index.js
   ;;
 stop)
   forever stop vaerhona-staging
   ;;
 restart)
   PORT=8001 NODE_ENV=production AWS_PROFILE=staging MONGO_URL=mongodb://localhost:27017/vaerhona-staging forever restart --sourceDir=$APP_DIR -w --watchDirectory $APP_DIR/dist --silent --uid "vaerhona-staging" --append index.js
   ;;
 *)
   echo "Usage: vaerhonastaging {start|stop|restart}" >&2
   exit 3
   ;;
esac